# dot_config/systemd/user/myBackup.timer
[Unit]
Description=Timer for Backup utility

[Timer]
OnCalendar=Sun *-*-* *:00/30:00

[Install]
WantedBy=timers.target
