# dot_config/systemd/user/swayidle.service
[Unit]
Description=SwayIdle Service
# Start the agent after the graphical session is ready (optional, good practice)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/sbin/swayidle -w \
timeout 600 'qs -c noctalia-shell ipc call lockScreen lock' \
timeout 630 'niri msg action power-off-monitors' \
timeout 1200 'qs -c noctalia-shell ipc call sessionMenu lockAndSuspend'
# Ensure the process is restarted if it fails
Restart=on-failure
# Give the agent time to register with the D-Bus session
TimeoutSec=30

[Install]
WantedBy=graphical-session.target
